# target/linux/pxa/zipitz2/base-files/etc/inittab
::sysinit:/usr/sbin/loadfont < /etc/console-setup/Uni2-Terminus12x6.psf
::sysinit:/usr/bin/loadkeys /etc/console-setup/keymap-2.6.map
::sysinit:/etc/init.d/rcS S boot
::shutdown:/etc/init.d/rcS K shutdown
ttyS0::askfirst:/bin/ash --login
tty1::respawn:/usr/bin/z2-shell
tty2::askfirst:/bin/ash --login
