# target/linux/pxa/zipitz2/base-files/etc/profile.d/zipit-console.sh
export TERM=xterm-color
export SDL_NOMOUSE=1
