# target/linux/pxa/zipitz2/base-files/etc/init.d/zipit-console
#!/bin/sh /etc/rc.common
# Copyright (C) 2008-2010 OpenWrt.org

START=18

start() {
    # set up some details to make the Z2 console legible    
    /usr/sbin/loadfont < /etc/console-setup/Uni2-Terminus12x6.psf
    /sbin/loadkmap < /etc/z2.kmap
}
